# deploy/aws/systemd/workbench-pipeline-health-weekly.timer
[Unit]
Description=Pipeline health WEEKLY: Monday 11:15 ET, after ALL four rebalances (10:00-10:40)
[Timer]
OnCalendar=Mon 11:15
AccuracySec=1min
Persistent=false
[Install]
WantedBy=timers.target
