# deploy/aws/systemd/workbench-pipeline-health-weekly.service
[Unit]
Description=Pipeline health checklist (WEEKLY): data freshness + rebalance status -> DB + SNS
After=docker.service
Wants=docker.service
[Service]
Type=oneshot
Environment=AWS_REGION=us-east-1
WorkingDirectory=/opt/workbench/app
ExecStart=/bin/bash /opt/workbench/app/deploy/aws/pipeline-health.sh WEEKLY
SuccessExitStatus=0 1 2
